# initfs/etc/inittab
ttyS0::sysinit:/etc/greeting
ttyS0::respawn:/bin/sh
